# scripts/systemd/opdd.service
[Unit]
Description=OPD Panel Daemon
After=network.target opd-agent.service
Requires=opd-agent.service

[Service]
Type=simple
EnvironmentFile=/etc/opd/.env
ExecStart=/usr/local/bin/opdd
Restart=on-failure
RestartSec=5s
User=root
Group=opd
RuntimeDirectory=opd
RuntimeDirectoryMode=0750

[Install]
WantedBy=multi-user.target
